# CP_Source/Tools_Drivers/Magnus_Screen_Magnifier-2/build/igel-magnus-cp-apparmor-reload.service
[Unit]
Description=Magnus Screen Magnifier custom partition apparmor reload
After=sysinit.target
After=apparmor.service
Before=igel-sessions-config.target
AssertPathIsReadWrite=/sys/kernel/security/apparmor/.load
ConditionSecurity=apparmor
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/config/bin/magnus_cp_apparmor_reload
